# src/etc/cron.d/sa-learn
*/5 * * * * debian-spamd /srv/sa-learn.sh ${MAIL_DIR}
