# src/etc/cron.d/sa-learn
0 * * * * root /srv/sa-learn.sh ${MAIL_DIR}
