# etc/cron.d/learn-spam
# Learn SPAM (bad mails)
0 4 * * * sa-learn --dbpath /var/lib/spamassassin/.spamassassin --spam ${MAIL_DIR}/*/*/*/.Spam/cur/*
