# etc/cron.d/learn-ham
0 3 * * * debian-spamd sa-learn --ham ${MAIL_DIR}/*/*/*/cur/*
