# etc/cron.d/learn-spam
0 4 * * * debian-spamd sa-learn --spam ${MAIL_DIR}/*/*/*/.Spam/cur/*
