# etc/cron.d/learn-ham
# Learn HAM (good mails)
0 3 * * * sa-learn --dbpath /var/lib/spamassassin/.spamassassin --ham ${MAIL_DIR}/*/*/*/cur/*
